# sensor_modules/services/boat.target
[Unit]
Description=Main Boat target for services
Requires=default.target
After=default.target
